# ansible/roles/ispstats/templates/ispstats.timer
[Unit]
Description=ISP Stats

[Timer]
OnActiveSec=10s
OnUnitActiveSec=1h

[Install]
WantedBy=basic.target
