# ansible/roles/ispstats/templates/ispstats.service.j2
[Unit]
Description=ISP Stat Script

[Service]
User={{ ispstats_user }}
Group={{ ispstats_user }}
ExecStart={{ ispstats_path }}/stats.py
