# ansible/roles/ispstats/templates/ispstats.service
[Unit]
Description=ISP Stat Script

[Service]
User=stats
Group=stats
ExecStart=/opt/ispstats/stats.py
