# ansible/roles/ispstats/templates/ispstats.timer.j2
[Unit]
Description=isp-stats Timer

[Timer]
OnCalendar={{ run_interval }}
Persistent=true

[Install]
WantedBy=basic.target
