# ansible/roles/ispstats/templates/ispstats.service.j2
[Unit]
Description=isp-stats Service

[Service]
User={{ ispstats_user }}
Group={{ ispstats_user }}
Nice=-20
SyslogIdentifier=isp-stats
ExecStart={{ ispstats_path }}/stats.py
