# ansible/roles/ispstats/templates/ispstats.timer.j2
[Unit]
Description=ISP Stats

[Timer]
OnActiveSec=10s
OnUnitActiveSec={{ run_interval }}

[Install]
WantedBy=basic.target
